# src/delete_files/delete_fonts.rules
KERNEL=="loop0", RUN+="/root/delete_fonts.sh"
